# systemd/coffee-reminder.timer
[Unit]
Description=Coffee Reminder timer (daily at 19:40)

[Timer]
OnCalendar=*-*-* 19:40:00
Persistent=true
Unit=coffee-reminder.service

[Install]
WantedBy=timers.target
