# systemd/coffee-reminder.service
[Unit]
Description=Coffee Reminder (with kernel beep)

[Service]
Type=oneshot
# Try to load kernel module and configure it; ignore errors in user context
ExecStart=%h/Coffee_remainder/coffee_reminder.sh

[Install]
WantedBy=default.target
